# test-firmware/_openwrt-25.12.5-ath79-generic-8dev_lima-squashfs-sysupgrade.bin.extracted/squashfs-root-0/lib/apk/packages/luci-base.list
/etc/config/luci
/etc/init.d/ucitrack
/etc/luci-uploads/.placeholder
/lib/upgrade/keep.d/luci-base
/usr/lib/ucode/luci/core.so
/usr/share/acl.d/luci-base.json
/usr/share/luci/menu.d/luci-base.json
/usr/share/rpcd/acl.d/luci-base.json
/usr/share/rpcd/ucode/luci
/usr/share/ucode/luci/controller/admin/index.uc
/usr/share/ucode/luci/controller/admin/uci.uc
/usr/share/ucode/luci/dispatcher.uc
/usr/share/ucode/luci/http.uc
/usr/share/ucode/luci/runtime.uc
/usr/share/ucode/luci/sys.uc
/usr/share/ucode/luci/template/csrftoken.ut
/usr/share/ucode/luci/template/error404.ut
/usr/share/ucode/luci/template/error500.ut
/usr/share/ucode/luci/template/footer.ut
/usr/share/ucode/luci/template/header.ut
/usr/share/ucode/luci/template/sysauth.ut
/usr/share/ucode/luci/template/view.ut
/usr/share/ucode/luci/uhttpd.uc
/usr/share/ucode/luci/version.uc
/usr/share/ucode/luci/zoneinfo.uc
/www/cgi-bin/luci
/www/index.html
/www/luci-static/resources/cbi.js
/www/luci-static/resources/cbi/file.svg
/www/luci-static/resources/cbi/folder.svg
/www/luci-static/resources/cbi/link.svg
/www/luci-static/resources/firewall.js
/www/luci-static/resources/form.js
/www/luci-static/resources/fs.js
/www/luci-static/resources/icons/alias.svg
/www/luci-static/resources/icons/alias_disabled.svg
/www/luci-static/resources/icons/bridge.svg
/www/luci-static/resources/icons/bridge_disabled.svg
/www/luci-static/resources/icons/ethernet.svg
/www/luci-static/resources/icons/ethernet_disabled.svg
/www/luci-static/resources/icons/loading.svg
/www/luci-static/resources/icons/port_down.svg
/www/luci-static/resources/icons/port_pse_down.svg
/www/luci-static/resources/icons/port_pse_up.svg
/www/luci-static/resources/icons/port_up.svg
/www/luci-static/resources/icons/signal-000-000.svg
/www/luci-static/resources/icons/signal-000-025.svg
/www/luci-static/resources/icons/signal-025-050.svg
/www/luci-static/resources/icons/signal-050-075.svg
/www/luci-static/resources/icons/signal-075-100.svg
/www/luci-static/resources/icons/signal-none.svg
/www/luci-static/resources/icons/switch.svg
/www/luci-static/resources/icons/switch_disabled.svg
/www/luci-static/resources/icons/tunnel.svg
/www/luci-static/resources/icons/tunnel_disabled.svg
/www/luci-static/resources/icons/vlan.svg
/www/luci-static/resources/icons/vlan_disabled.svg
/www/luci-static/resources/icons/vrf.svg
/www/luci-static/resources/icons/vrf_disabled.svg
/www/luci-static/resources/icons/wifi.svg
/www/luci-static/resources/icons/wifi_disabled.svg
/www/luci-static/resources/icons/wireguard.svg
/www/luci-static/resources/icons/wireguard_disabled.svg
/www/luci-static/resources/luci.js
/www/luci-static/resources/network.js
/www/luci-static/resources/protocol/dhcp.js
/www/luci-static/resources/protocol/none.js
/www/luci-static/resources/protocol/static.js
/www/luci-static/resources/rpc.js
/www/luci-static/resources/tools/prng.js
/www/luci-static/resources/tools/views.js
/www/luci-static/resources/tools/widgets.js
/www/luci-static/resources/uci.js
/www/luci-static/resources/ui.js
/www/luci-static/resources/validation.js
/www/luci-static/resources/xhr.js
