# test-firmware/_openwrt-25.12.5-ath79-generic-8dev_carambola3-squashfs-sysupgrade.bin.extracted/squashfs-root/lib/apk/packages/hostapd-common.list
/etc/capabilities/wpad.json
/etc/init.d/wpad
/lib/netifd/dhcp-get-server.sh
/usr/share/acl.d/wpad_acl.json
